# container/ubooquity/build/supervisord.ini
[supervisord]
nodaemon=true

[program:ubooquity]
command=java -jar /opt/ubooquity/Ubooquity.jar -webadmin -headless -workdir /home/ubooquity
redirect_stderr=true
user=ubooquity

#[program:httpd]
#command=/usr/sbin/httpd -D FOREGROUND -c "ErrorLog /dev/stderr" -c "CustomLog /dev/stderr common"
#redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx
redirect_stderr=true