# container/backuppc/build/supervisord.ini
[supervisord]
nodaemon=true

[program:sshd]
command=/sbin/execute_sshd.sh
redirect_stderr=true

[program:httpd]
command=/usr/sbin/httpd -D FOREGROUND -c "ErrorLog /dev/stderr" -c "CustomLog /dev/stderr common"
redirect_stderr=true

[program:backuppc]
command=/usr/share/BackupPC/bin/BackupPC
redirect_stderr=true
user={{backuppc_user}}