# container/library/files/library.service.j2
[Unit]
Description=library volume Service
After=etcd.service docker.service
Requires=etcd.service docker.service

[Service]
EnvironmentFile=/etc/sysconfig/dnas
ExecStartPre=-/usr/bin/docker kill library
ExecStartPre=-/usr/bin/docker pull sebastianhutter/library
ExecStartPre=-/usr/bin/docker create --name library \
  --volume /Library \
  --volume /Transfer \
  sebastianhutter/library \
  library
ExecStart=/usr/bin/docker start library
ExecStop=/usr/bin/docker stop library
Restart=on-failure
# up the timeout values
RestartSec=2s
TimeoutStartSec=10s


[Install]
WantedBy=multi-user.target


