# container/library/files/library.service.j2
[Unit]
Description=library volume Service
After=etcd.service docker.service
Requires=etcd.service docker.service

[Service]
EnvironmentFile=/etc/sysconfig/dnas
ExecStartPre=-/usr/bin/docker kill library
ExecStartPre=/usr/bin/docker pull sebastianhutter/library
ExecStart=/usr/bin/docker run --name library \
  --volume /Library \
  sebastianhutter/library \
  library
ExecStop=/usr/bin/docker stop library
Restart=never

[Install]
WantedBy=multi-user.target


